# resources/supervisor/mongod.conf
[program:mongod]
command=python /opt/mongod/pod/pod.py
